# deploy/systemd/life-reminders.service
[Unit]
Description=Life System reminder sender
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/life-system
EnvironmentFile=/etc/life-system/life-system.env
ExecStart=/opt/life-system/scripts/run_reminders.sh

